# archlinux/var/lib/pacman/local/fake-hwclock-0.3-1/files
%FILES%
usr/
usr/lib/
usr/lib/systemd/
usr/lib/systemd/scripts/
usr/lib/systemd/scripts/fake-hwclock.sh
usr/lib/systemd/system/
usr/lib/systemd/system/fake-hwclock-save.service
usr/lib/systemd/system/fake-hwclock-save.timer
usr/lib/systemd/system/fake-hwclock.service

