# archlinux/usr/lib/systemd/system/rce_dma.service
[Unit]
Description=RCE DMA Driver
Requires=mnt-host.mount
After=mnt-host.mount

[Service]
Type=forking
ExecStart=/mnt/host/heavy-photon-daq/software/rceScripts/rce_dma.sh

[Install]
WantedBy=multi-user.target

