# archlinux/usr/lib/systemd/system/rce_start.service
[Unit]
Description=RCE Startup Script
After=network-online.target

[Service]
#Type=forking
Type=oneshot
ExecStart=/mnt/host/heavy-photon-daq/software/rceScripts/rce_start.sh

[Install]
WantedBy=multi-user.target

