# archlinux/usr/lib/systemd/system/rce_halt.timer
[Unit]
Description=Run to halt RCEs.

[Timer]
OnActiveSec=1s
Unit=rce_halt.service

[Install]
WantedBy=multi-user.target

