# archlinux/usr/lib/systemd/system/rce_start.service
[Unit]
Description=RCE Startup Script
After=network-online.target

[Service]
Type=forking
ExecStart=/mnt/host/daq/rceData/rce_start.sh

[Install]
WantedBy=multi-user.target

